# dpkg-build/DEBIAN/prerm
#!/bin/sh
systemctl disable dotnet-veml7700